# egpu.service
[Unit]
Description=EGPU Service

[Service]
ExecStart=egpu-switcher switch auto

[Install]
WantedBy=multi-user.target
