# egpu.service
[Unit]
Description=EGPU Service

[Service]
Type=oneshot
ExecStart=/usr/bin/egpu-switcher switch auto

[Install]
WantedBy=multi-user.target
